# build/scripts/files/etc/uci-defaults/aaa.sh
#!/bin/sh

cp -Rf /etc/openclash/custom/ /root/openclash_customaaa

exit 0
